# com.heartflow.enhanced-upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heartflow.enhanced-upgrade</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/apple/mark-heartflow-skill-new/scripts/hourly-enhanced-upgrade-v2.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1440</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/apple/mark-heartflow-skill-new/logs/enhanced-cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/apple/mark-heartflow-skill-new/logs/enhanced-cron-error.log</string>
</dict>
</plist>
