# com.heartflow.auto-upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heartflow.auto-upgrade</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/apple/mark-heartflow-skill/scripts/auto-upgrade.js</string>
    </array>
    <StartCalendarInterval>
      <dict><key>Minute</key><integer>0</integer><key>Hour</key><integer>*</integer></dict>
      <dict><key>Minute</key><integer>15</integer><key>Hour</key><integer>*</integer></dict>
      <dict><key>Minute</key><integer>30</integer><key>Hour</key><integer>*</integer></dict>
      <dict><key>Minute</key><integer>45</integer><key>Hour</key><integer>*</integer></dict>
    </StartCalendarInterval>
    <key>StandardOutPath</key>
    <string>/Users/apple/mark-heartflow-skill/logs/auto-upgrade.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/apple/mark-heartflow-skill/logs/auto-upgrade.error.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
