# .opencode/heartflow.cron
# HeartFlow 定时任务配置
# 每 30 分钟运行一次
# 使用 MiniMax M2.5 模型

# 格式：分 时 日 月 星期 命令

# 每 30 分钟运行自动运行器 (MiniMax M2.5)
*/30 * * * * cd ~/.jvs/.openclaw/workspace/mark-heartflow-skill && MODEL=MiniMax-M2.5 node .opencode/auto-runner.js >> .opencode/cron.log 2>&1

# 每天凌晨 2 点生成日报 (MiniMax M2.5)
0 2 * * * cd ~/.jvs/.openclaw/workspace/mark-heartflow-skill && MODEL=MiniMax-M2.5 node -e "const {generateReport}=require('./.opencode/memory/report-generator');const {initializeState}=require('./.opencode/memory/init');const s=initializeState('./.opencode/memory/heartflow_state.json');console.log(generateReport(s));" >> .opencode/daily-report.log 2>&1

# 每周日清理旧日志
0 3 * * 0 find ~/.jvs/.openclaw/workspace/mark-heartflow-skill/.opencode -name "*.log" -mtime +7 -delete 2>/dev/null
