# internal/com.heartflow.auto-upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heartflow.auto-upgrade</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/apple/mark-heartflow-skill/auto-hourly-upgrade.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1380</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/apple/mark-heartflow-skill/logs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/apple/mark-heartflow-skill/logs/launchd-error.log</string>
</dict>
</plist>